# components/desktop-env/files/etc/debcomp/userconf/desktop/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=/usr/local/share/dasc/libexec/desktop/start-conky
