# components/desktop-env-mail-client/files/etc/debcomp/userconf/desktop/skel/.config/autostart/claws-mail.desktop
[Desktop Entry]
Type=Application
Name=Claws Mail
Icon=claws-mail
Exec=claws-mail
NotShowIn=OPENBOX;
